# etc/systemd/system/firewalld.service.d/restart-docker.conf
[Service]
ExecStartPost=/usr/bin/systemctl restart --no-block docker
ExecReload=
ExecReload=/bin/kill -HUP $MAINPID
ExecReload=/usr/bin/systemctl restart --no-block docker
