# etc/systemd/system/docker.service.d/override.conf
[Service]
ExecStartPre=/bin/sh -c "until /usr/bin/firewall-cmd --state >/dev/null 2>&1; do sleep 1; done"