# etc/systemd/system/nginx.service.d/override.conf
[Unit]
StartLimitIntervalSec=0

[Service]
Restart=on-failure
RestartSec=20s