# etc/systemd/system/nftables-mailcow-net.service
[Unit]
Description=Mailcow DNAT + restricted egress NAT (preserve real client IPs)
After=firewalld.service docker.service
Wants=firewalld.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/sbin/nft-mailcow-net.sh

[Install]
WantedBy=multi-user.target
